# brain/scripts/discord_inbox_bot/launchd/com.sasaki.discord-inbox-bot.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sasaki.discord-inbox-bot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>VAULT_PATH/brain/scripts/discord_inbox_bot/bot.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>VAULT_PATH/brain/scripts/discord_inbox_bot</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>VAULT_PATH/brain/scripts/discord_inbox_bot/discord-inbox-bot.log</string>
  <key>StandardErrorPath</key>
  <string>VAULT_PATH/brain/scripts/discord_inbox_bot/discord-inbox-bot.err</string>
</dict>
</plist>
